# NIMROD_K/2.0.0/etc/profile.d/nimrod_g_modules.sh
#!/bin/bash
echo -n "Checking NIMROD G 'modules' requirements..."

# set paths
export NIMRODJAVAAPI_INSTALL=/usr/local/kepler/2.4/lib/nimrodk-2.4.0
export KEPLER=/usr/local/kepler/2.4/lib

echo -e "Finished.\n"
